# emoncms/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/bin/bash -c "exec mysqld_safe"

[program:redis]
command=/bin/bash -c "exec redis-server"

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"
